# deploy/stockbot.service
[Unit]
Description=StockBot Trading API
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/stockbot
ExecStart=/usr/bin/docker compose up --build
ExecStop=/usr/bin/docker compose down
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stockbot

[Install]
WantedBy=multi-user.target
